# scripts/setup/telegram-bot.service
[Unit]
# 示例 systemd 单元，请按部署环境调整 User / WorkingDirectory / EnvironmentFile
Description=FateCat Telegram Bot
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=fatecat
WorkingDirectory=/srv/fatecat/modules/telegram
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/etc/fatecat/fatecat.env
ExecStart=/bin/bash -lc 'source /srv/fatecat/.venv/bin/activate && python start.py bot'
Restart=always
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
